# OpenWrt/car_controller.service
[Unit]
Description=Smart Car Controller Service
After=network.target serial-getty@ttyUSB0.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/car_controller
ExecStart=/usr/bin/python3 /opt/car_controller/socket_server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
